# pkg/etc/profile.d/kamilscripts.sh

if [ -e /usr/lib/kamilscripts/etc/profile ] &&
	[ -r /usr/lib/kamilscripts/etc/profile ]; then
	. /usr/lib/kamilscripts/etc/profile /usr/lib/kamilscripts/etc
fi
