# systemd/system/_snapper_hourly.timer
[Unit]
Description=Run _snapper_hourly hourly

[Timer]
OnCalendar=Hourly
AccuracySec=5m
Persistent=true

[Install]
WantedBy=multi-user.target
