# systemd/system/mnt-perunshare.automount
[Unit]
Description=Automount /mnt/perunshare

[Automount]
Where=/mnt/perunshare

[Install]
WantedBy=multi-user.target
