# systemd/system/qqremember.timer
[Unit]
Description=Run qqremember Daily

[Timer]
OnCalendar=*-*-* 13:00:00
AccuracySec=1h
Persistent=true

[Install]
WantedBy=multi-user.target
