# stow/noinstall/.local/share/profile.d/kamilscripts.sh
../../../../../etc/profile.d/kamilscripts.sh