# resources/usr/lib/systemd/system/gitsbackup.sh.timer
[Unit]
Description = Daily gitsbackup.sh run

[Timer]
OnCalendar=Daily
AccuracySec=1d
Persistent=true
