# resources/usr/lib/systemd/system/pacman-optimize.timer
[Unit]
Description=Run pacman-optimize monthly

[Timer]
OnCalendar=monthly
AccuracySec=1d
Persistent=true

[Install]
WantedBy=multi-user.target
