# resources/usr/lib/systemd/system/qqbackuphistfile.sh@.timer
[Unit]
Description=Run backuphistfile.sh@%I Daily

[Timer]
OnCalendar=Daily
AccuracySec=1d
Persistent=true

[Install]
WantedBy=multi-user.target
