# resources/usr/lib/systemd/system/qqrdiffbackup.sh@.timer
[Unit]
Description=Run qqrdiffbackup.sh@%I Daily

[Timer]
OnCalendar=Daily
AccuracySec=1d
Persistent=true

[Install]
WantedBy=multi-user.target
