# resources/usr/lib/systemd/system/qqremember.timer
[Unit]
Description=Run qqremember Daily

[Timer]
OnCalendar=*-*-* 06:00:00
AccuracySec=1d
Persistent=true

[Install]
WantedBy=multi-user.target
