# etc/profile.d/awkpath.sh
#!/bin/sh

# https://www.gnu.org/software/gawk/manual/html_node/Include-Files.html
AWKPATH="$KCDIR/awk${AWKPATH:+:$AWKPATH}"
export AWKPATH

