# etc/profile.d/docker.sh
#!/bin/sh
# Remove docker ad...
export BUILDKIT_PROGRESS=plain
export DOCKER_SCAN_SUGGEST=false
