# etc/profile.d/docker.sh
#!/bin/sh
# Remove docker ad...
export DOCKER_SCAN_SUGGEST=false
